# services/ssh-guard-init.service
[Unit]
Description=SSH Guard initialisation – master key generation
After=local-fs.target
Before=ssh-guard.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ssh-guard --genkey
User=root
Group=root
ProtectSystem=strict
ReadWritePaths=/etc/ssh-guard
NoNewPrivileges=yes
PrivateNetwork=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ssh-guard-init

[Install]
WantedBy=multi-user.target