# services/ssh-guard.service
[Unit]
Description=SSH directory access guard (fanotify)
Documentation=man:fanotify(7)
After=local-fs.target ssh-guard-init.service
Before=sshd.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/ssh-guard
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2s

User=root
Group=root

PrivateTmp=yes
PrivateNetwork=yes
ProtectSystem=yes
ReadWritePaths=/run /etc/ssh-guard
NoNewPrivileges=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ssh-guard

[Install]
WantedBy=multi-user.target